# deport/wfw/wildbook_supervisor.conf
[program:wildbook]
command=uwsgi /home/ubuntu/virtualenv/wildbook/wildbook/wildteam/deport/wfw/wildteam_uwsgi_wfw.ini -d log/uwsgi.log